# systemd/town-os-rolodex-config.service
# vim: ft=systemd
[Unit]
Description=Refresh DNS (bootstrap resolver + rolodex binds) from the current DHCP lease
After=town-os-system--rolodex.service
RequiresMountsFor=/town-os

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/lib/town-os/scripts/rolodex-reload.sh
