# systemd/town-os-systemcontroller.service
# vim: ft=systemd
[Unit]
Description=Town OS System Controller
DefaultDependencies=yes
Requires=network-online.target
Wants=town-os-system--rolodex.service
Wants=town-os-podman-api.service
After=network-online.target town-os-system--rolodex.service town-os-podman-api.service
RequiresMountsFor=/town-os
StartLimitIntervalSec=0

[Service]
Type=exec
ExecStartPre=/bin/ln -sf /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
ExecStartPre=/usr/bin/mkdir -p /town-os/data/db /town-os/data/repos /var/lib/containers /run/town-os
ExecStartPre=-/usr/bin/podman rm -f systemcontroller
# --pull=missing, NOT --pull=always: on a crash/restart (e.g. the systemd
# watchdog killing us after a suspend/resume freeze) --pull=always blocks the
# restart on re-pulling our own image from quay.io; if the host network isn't
# back yet the pull fails and the controller can't restart, so the box loses
# management until the internet returns. --pull=missing restarts from the local
# image with no network dependency. Updates are still applied deliberately: the
# UI's POST /system-services/refresh explicitly `podman pull`s the controller
# image (it's in collectSystemServices) before scheduling the self-restart, so
# the freshly-pulled image is already local when --pull=missing runs it.
ExecStart=/usr/bin/podman run --pull=missing \
    --name systemcontroller \
    --privileged \
    --security-opt label=disable \
    --net host \
    --pid host \
    -v /etc/systemd:/etc/systemd \
    -v /usr/lib/systemd:/usr/lib/systemd \
    -v /run/systemd:/run/systemd \
    -v /run/dbus:/run/dbus \
    -v /run/podman/podman.sock:/run/podman/podman.sock \
    -v /var/lib/containers:/var/lib/containers:shared \
    -v /dev:/dev \
    -v /etc/machine-id:/etc/machine-id:ro \
    -v /town-os:/town-os:shared \
    -v /var/log:/var/log \
    -v /run/town-os:/run/town-os \
    -e TOWN_OS_TAG=@TOWN_OS_TAG@ \
    quay.io/town/town:rc.latest \
    /systemcontroller \
    -repo-dir /town-os/data/repos \
    -db /town-os/data/db/system.db \
    -btrfs /town-os \
    @PACKAGE_DNS@
ExecStopPost=/bin/ln -sf /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
Restart=always
RestartSec=5
RestartForceExitStatus=143 SIGTERM

[Install]
WantedBy=default.target
