# systemd/town-os-network-diag.service
# vim: ft=systemd
[Unit]
Description=Town OS Network Diagnostics
After=town-os-systemcontroller.service

[Service]
# A long-lived recorder, not a timer-driven oneshot: it emits one journal line
# at boot and is then silent, instead of a Starting/Deactivated/Finished triple
# every interval. It appends to its log only when the network state changes.
Type=simple
ExecStart=/bin/sh /usr/lib/town-os/scripts/network-diag.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
