# systemd/town-os-system--rolodex.service
# vim: ft=systemd
[Unit]
Description=Town OS System Service: Rolodex DNS
DefaultDependencies=yes
Requires=network-online.target
After=network-online.target
RequiresMountsFor=/town-os
Before=town-os-systemcontroller.service
StartLimitIntervalSec=0

[Service]
Type=exec
# Generate rolodex.yml (dynamic bind list; no forwarders — rolodex recurses from
# the root servers). See the script for the full rationale.
ExecStartPre=/bin/sh /usr/lib/town-os/scripts/rolodex-config.sh
ExecStartPre=-/usr/bin/podman container cleanup town-os-system--rolodex
ExecStartPre=-/usr/bin/podman stop -t 10 town-os-system--rolodex
ExecStartPre=-/usr/bin/podman rm -f town-os-system--rolodex
ExecStart=/usr/bin/podman run --pull=always --replace --name town-os-system--rolodex \
    --net host \
    -v /town-os/rolodex:/data \
    quay.io/town/rolodex:rc.latest \
    /usr/local/bin/rolodex-dns --config /data/rolodex.yml
ExecStop=/usr/bin/podman stop -t 10 town-os-system--rolodex
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
