# systemd/town-os-system--rolodex.service
# vim: ft=systemd
[Unit]
Description=Town OS System Service: Rolodex DNS
DefaultDependencies=yes
Requires=network-online.target
After=network-online.target
RequiresMountsFor=/town-os
Before=town-os-systemcontroller.service
StartLimitIntervalSec=0

[Service]
Type=exec
ExecStartPre=/bin/mkdir -p /town-os/rolodex
ExecStartPre=/bin/sh -c '\
  EXTERNAL_IP=$(ip -4 route get 1.1.1.1 2>/dev/null | awk "/src/{for(i=1;i<=NF;i++) if($i==\"src\") print $(i+1)}"); \
  printf "database_path: /data/rolodex.db\ndns:\n  udp_bind:\n    - \"127.0.0.2:53\"\n    - \"${EXTERNAL_IP:-0.0.0.0}:53\"\n  tcp_bind:\n    - \"127.0.0.2:53\"\n    - \"${EXTERNAL_IP:-0.0.0.0}:53\"\ngrpc:\n  tcp_bind: \"\"\n  unix_socket: /data/rolodex.sock\n  shared_secret: \"\"\nforwarders:\n  - \"1.1.1.1:53\"\n  - \"8.8.8.8:53\"\nrbl:\n  enabled: true\n  providers: []\n" > /town-os/rolodex/rolodex.yml'
ExecStartPre=-/usr/bin/podman container cleanup town-os-system--rolodex
ExecStartPre=-/usr/bin/podman stop -t 10 town-os-system--rolodex
ExecStartPre=-/usr/bin/podman rm -f town-os-system--rolodex
ExecStart=/usr/bin/podman run --pull=missing --replace --name town-os-system--rolodex \
    --net host \
    -v /town-os/rolodex:/data \
    quay.io/town/rolodex:rc.latest \
    /usr/local/bin/rolodex-dns --config /data/rolodex.yml
ExecStop=/usr/bin/podman stop -t 10 town-os-system--rolodex
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
