# systemd/avahi-daemon.service.d/network.conf
[Unit]
After=network-online.target town-os-systemcontroller.service
Wants=network-online.target town-os-systemcontroller.service
